# random-dwall.timer
[Unit]
Description=Run random-dwall hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
